# files_lumi_domoticz/etc/uci-defaults/80_add_domoticz_dialout
. /lib/functions.sh

group_add_user "dialout" "domoticz"

exit 0
