# files_lumi_domoticz/etc/uci-defaults/80_domoticz_userdata
sed -i -e "s:option userdata .*:option userdata '/etc/domoticz/':"     /etc/config/domoticz
chown -R domoticz:domoticz "/etc/domoticz"

exit 0